# mochad.conf
[program:start]
autorestart = false
startsecs = 0
user = root
command = /usr/local/bin/mochad