# mochad.conf
[program:start]
autorestart = false
startsecs = 10
user = root
command = /usr/local/bin/mochad -d
stdout_logfile=/var/log/mochad.log
stderr_logfile=/var/log/mochad.err
