# mochad.conf
[program:start]
autorestart = false
startsecs = 10
user = root
command = /usr/local/bin/mochad -d
redirect_stderr = true
stdout_logfile=/var/log/mochad.log
