# services/mochad/run
#!/bin/bash

tail -F /var/log/mochad.log
